# .config/.crontab
0 */3 * * * /usr/bin/crontab -l > /Users/ce/.config/.crontab
0 */3 * * * /bin/zsh /Users/ce/.scripts/conda.sh
0 */3 * * * /opt/homebrew/bin/brew bundle dump --force --file=/Users/ce/.brew


